# deploy/systemd/wildboar-worker-compliance.service
[Unit]
Description=WildBoar Worker - Compliance Checker
After=network.target

[Service]
Type=simple
User=wildboar
Group=wildboar
WorkingDirectory=/opt/wildboar/current
EnvironmentFile=/opt/wildboar/shared/.env
ExecStart=/opt/wildboar/.venv/bin/python -m workers.bsc_compliance_checker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target