# deploy/systemd/wildboar-worker-deposit-listener.service
[Unit]
Description=WildBoar Worker - Deposit Listener
After=network.target

[Service]
Type=simple
User=wildboar
Group=wildboar
WorkingDirectory=/opt/wildboar
EnvironmentFile=/etc/wildboar/wildboar.env
ExecStart=/opt/wildboar/.venv/bin/python -m workers.bsc_usdt_deposit_listener
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target