# deploy/systemd/wildboar-worker-telegram-watchdog.service
[Unit]
Description=WildBoar Worker - Telegram Fee Wallet Watchdog
After=network.target

[Service]
Type=simple
User=wildboar
Group=wildboar
WorkingDirectory=/opt/wildboar/current
EnvironmentFile=/opt/wildboar/shared/.env
ExecStart=/opt/wildboar/.venv/bin/python -m workers.telegram_fee_wallet_watchdog
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target