# deploy/systemd/wildboar-worker-balance-updater.service
[Unit]
Description=WildBoar Worker - USDT Balance Updater
After=network.target

[Service]
Type=simple
User=wildboar
Group=wildboar
WorkingDirectory=/opt/wildboar
EnvironmentFile=/etc/wildboar/wildboar.env
ExecStart=/opt/wildboar/.venv/bin/python -m workers.bsc_usdt_balance_updater
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target